# efixo-base-files/files/nb5-refurbishing/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
::respawn:/etc/init.d/telnet start
ttyS0::askfirst:-/bin/login
