# efixo-base-files/files/rescue/etc/init.d/autoconf
#!/bin/sh /etc/rc.common

nohup /usr/sbin/autoconf.sh > /dev/console
