# efixo-base-files/files/rescue/etc/init.d/switch
#!/bin/sh /etc/rc.common

boot() {
	case "`hostname`" in
		nb6)
			rtk-switch init
			;;
	esac
}
