# efixo-lan-topology/files/lan-topology.init
#!/bin/sh /etc/rc.common

boot() {
	ip neigh flush all
	lan-topology
}
